# root/etc/services.d/pwndrop/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc \
		/app/pwndrop/pwndrop \
			-debug \
			-no-autocert \
			-no-dns \
			-config /defaults/pwndrop.ini
